# files/system/etc/skel/.config/autostart/org.kde.krdc.desktop
[Desktop Entry]
Exec=krdc rdp://GH_SEC_RDP_ENDPOINT
Icon=krdc
Name[fr_BE]=KRDC
Name=KRDC
StartupWMClass=krdc
Terminal=false
Type=Application
X-DocPath=krdc/index.html
