# files/system/etc/systemd/system/wg-recreate.service
[Unit]
Description=Recreate wireguard config on boot
Before=wg-quick@.service
ConditionPathExists=/usr/lib/env

[Service]
Type=oneshot
ExecStart=/usr/libexec/wg-recreate.sh

[Install]
WantedBy=wg-quick@wg0.service
