# files/system/etc/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
OnBootSec=5m
RandomizedDelaySec=1h
